# systemd/premierapi.service.d/override.conf
[Unit]
OnFailure=premierapi-startup-alert.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
EnvironmentFile=/etc/premierapi/premierapi.env
EnvironmentFile=/etc/premierapi/telegram-alerts.env
